# supervisord.conf
[supervisord]
nodaemon=true

[program:/loop/expressionevaluation.jar]
command=java -jar /loop/expressionevaluation.jar


[program:/loop/filepyramid.jar]
command=java -jar /loop/filepyramid.jar


